# build.service
[Unit]
Description=Startup build service

[Service]
Type=oneshot
ExecStart=/vagrant/setup.sh

[Install]
WantedBy=multi-user.target
