# crontab
1 0 * * * /vagrant/nightly.sh
0 0 1 * * /vagrant/monthly.sh

